# packaging/linux/deb/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        ln -sf /usr/lib/nattunnel/nattunneld /usr/bin/nattunneld
        ln -sf /usr/lib/nattunnel-gui/nattunnel-gui /usr/bin/nattunnel-gui

        # Service idles in Disconnected on boot (autoConnect=false). GUI triggers /connect.
        systemctl daemon-reload >/dev/null 2>&1 || true
        systemctl enable --now nattunnel.service >/dev/null 2>&1 || true
        ;;
esac

exit 0
